# task5/crontab
SHELL=/bin/bash
MAILTO=linux@otus.ru

0   *   *   *   *    cd /vagrant && ./nginx_stat.sh access.log
